# templates/default/initd.debian.erb
#! /bin/sh

### BEGIN INIT INFO
# Provides:          solr
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/stop Solr Server
# Description:       Start/stop Solr Server
### END INIT INFO

PATH=/sbin:/usr/sbin:/bin:/usr/bin
SOLR_DIR=<%= @solr_dir %>
PORT=<%= @port %>
DAEMON=<%= @solr_dir %>/bin/solr
DAEMON_ARGS="-p $PORT"
SCRIPTNAME=/etc/init.d/$NAME

# Exit gracefully if the package is not installed
[ -x "$DAEMON" ] || exit 0

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

#
# Function that starts the daemon/service
#
do_start()
{
  $DAEMON $DAEMON_ARGS
  sleep 2
}

#
# Function that stops the daemon/service
#
do_stop()
{
  $DAEMON stop
}

do_status()
{
  $DAEMON status
}

case "$1" in
  start)
  do_start
  ;;
  stop)
  do_stop
  ;;
  status)
  do_status
  ;;
  restart|force-reload)
  do_stop
  do_start
  ;;
  *)
  echo "Usage: $SCRIPTNAME {start|stop|status|restart|force-reload}" >&2
  exit 3
  ;;
esac

:
